# bin/crontab
#!/bin/bash

# 一次性 crontab 命令转发脚本。
# 该脚本保持和其他 bin 入口一致的调用方式，方便 Linux crontab 直接通过
# `scf/bin/crontab` 触发 `boot crontab`，而不需要显式拼接 boot 命令。

SCRIPT_DIR="$(cd "$(dirname "$0")" && pwd)"
BOOT_FILE="$SCRIPT_DIR/../boot"

PHP_BIN="${SCF_PHP_BIN:-}"
if [ -z "$PHP_BIN" ]; then
  PHP_BIN="$(command -v php 2>/dev/null || true)"
fi

if [ -z "$PHP_BIN" ]; then
  for candidate in /opt/homebrew/bin/php /usr/local/bin/php /usr/bin/php /opt/homebrew/opt/php*/bin/php /usr/local/opt/php*/bin/php; do
    if [ -x "$candidate" ]; then
      PHP_BIN="$candidate"
      break
    fi
  done
fi

if [ -z "$PHP_BIN" ] || [ ! -x "$PHP_BIN" ]; then
  echo "php binary not found, set SCF_PHP_BIN first" >&2
  exit 127
fi

# 解析是否由调用方显式指定了 pack 行为（-pack / -nopack），
# 以及当前命令的首个位置 action。
HAS_PACK_FLAG=0
ACTION=""
for arg in "$@"; do
  case "$arg" in
    -pack|-nopack)
      HAS_PACK_FLAG=1
      ;;
  esac
  case "$arg" in
    -*)
      ;;
    *)
      if [ -z "$ACTION" ]; then
        ACTION="$arg"
      fi
      ;;
  esac
done

if [ $# -eq 0 ]; then
  # 无参数时优先走源码命令实现，避免 pack 落后导致帮助信息缺失新 action。
  "$PHP_BIN" "$BOOT_FILE" crontab -nopack -h
  exit $?
fi

# `log` 只用于本地排查，默认走源码命令实现，确保能及时看到新增能力。
# 若调用方已显式传入 -pack/-nopack，则尊重其选择，不做注入。
INJECT_NO_PACK=0
if [ "$HAS_PACK_FLAG" -eq 0 ]; then
  if [ "$ACTION" = "log" ]; then
    INJECT_NO_PACK=1
  fi
fi

if [ "$INJECT_NO_PACK" -eq 1 ]; then
  "$PHP_BIN" "$BOOT_FILE" crontab -nopack "$@"
else
  "$PHP_BIN" "$BOOT_FILE" crontab "$@"
fi
