# bin/crontab
#!/bin/bash

# 一次性 crontab 命令转发脚本。
# 该脚本保持和其他 bin 入口一致的调用方式，方便 Linux crontab 直接通过
# `scf/bin/crontab` 触发 `boot crontab`，而不需要显式拼接 boot 命令。

SCRIPT_DIR="$(cd "$(dirname "$0")" && pwd)"
BOOT_FILE="$SCRIPT_DIR/../boot"

PHP_BIN="${SCF_PHP_BIN:-}"
if [ -z "$PHP_BIN" ]; then
  PHP_BIN="$(command -v php 2>/dev/null || true)"
fi

if [ -z "$PHP_BIN" ]; then
  for candidate in /opt/homebrew/bin/php /usr/local/bin/php /usr/bin/php /opt/homebrew/opt/php*/bin/php /usr/local/opt/php*/bin/php; do
    if [ -x "$candidate" ]; then
      PHP_BIN="$candidate"
      break
    fi
  done
fi

if [ -z "$PHP_BIN" ] || [ ! -x "$PHP_BIN" ]; then
  echo "php binary not found, set SCF_PHP_BIN first" >&2
  exit 127
fi

if [ $# -eq 0 ]; then
  "$PHP_BIN" "$BOOT_FILE" crontab -h
  exit $?
fi

"$PHP_BIN" "$BOOT_FILE" crontab "$@"
